# deploy/systemd/omnisec-nats.service
[Unit]
Description=OmniSec NATS (JetStream)
Documentation=https://github.com/manishbalayan/omnisec-v-0.1
After=network.target
PartOf=omnisec.target

[Service]
Type=simple
User=omnisec
Group=omnisec
ExecStart=/usr/local/libexec/omnisec/omnisec-nats-run
Restart=always
RestartSec=5
TimeoutStartSec=30
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ReadWritePaths=/var/lib/omnisec /var/log/omnisec

[Install]
WantedBy=omnisec.target
