# deploy/omnisec-daemon.service
[Unit]
Description=OmniSec Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/omnisec-daemon
Restart=always
RestartSec=10
# Environment variables - will be set by installer
Environment=DATABASE_URL=postgres://omnisec:omnisec@localhost:5432/omnisec
Environment=NATS_URL=nats://localhost:4222
Environment=OMNISEC_SAFE_MODE=0
Environment=OMNISEC_RECOMMENDATION_ONLY=0
Environment=RUST_LOG=info
# Watchdog - will be enabled if WATCHDOG_USEC is set
# WatchdogSec=30

[Install]
WantedBy=multi-user.target