# deploy/supervisord.conf
; =============================================================================
; OmniSec Supervisor Configuration
; =============================================================================
; This manages only the OmniSec application services.
; Infrastructure (PostgreSQL, NATS) is started by the entrypoint
; script BEFORE supervisor starts.
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/lib/omnisec/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/var/lib/omnisec/logs

; =============================================================================
; OmniSec Services
; =============================================================================

[program:omnisec-api]
command=omnisec-api
user=root
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/lib/omnisec/logs/api.log
stderr_logfile=/var/lib/omnisec/logs/api.log
environment=
    OMNISEC_API_KEY="%(ENV_OMNISEC_API_KEY)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    NATS_URL="%(ENV_NATS_URL)s",
    API_BIND="%(ENV_API_BIND)s"
priority=400

[program:omnisec-daemon]
command=omnisec-daemon
user=root
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/lib/omnisec/logs/daemon.log
stderr_logfile=/var/lib/omnisec/logs/daemon.log
environment=
    OMNISEC_API_KEY="%(ENV_OMNISEC_API_KEY)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    NATS_URL="%(ENV_NATS_URL)s",
    OMNISEC_SAFE_MODE="%(ENV_OMNISEC_SAFE_MODE)s",
    OMNISEC_RECOMMENDATION_ONLY="%(ENV_OMNISEC_RECOMMENDATION_ONLY)s",
    OMNISEC_VERBOSE="%(ENV_OMNISEC_VERBOSE)s",
    DAEMON_HEALTH_BIND="%(ENV_DAEMON_HEALTH_BIND)s"
priority=500

[program:omnisec-dashboard]
command=node /opt/omnisec-dashboard/server.js
user=root
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/lib/omnisec/logs/dashboard.log
stderr_logfile=/var/lib/omnisec/logs/dashboard.log
environment=
    PORT="%(ENV_DASHBOARD_PORT)s",
    HOSTNAME="0.0.0.0",
    NEXT_PUBLIC_API_URL="http://localhost:3002"
priority=600
